# rd98dx35xx_cn9131/service/rd98dx35xx-cn9131-init.service
[Unit]
Description=Platform rd98dx35xx-cn9131 Service
After=sysinit.target
Before=pmon.target

[Service]
ExecStart=/usr/local/bin/rd98dx35xx-cn9131-init.sh
StandardOutput=tty

[Install]
WantedBy=multi-user.target
